# sh/early-userspace/root/etc/inittab
::wait:/sbin/rc S
::ctrlaltdel:/sbin/shutdown -r
tty1::respawn:/bin/login -f root

